# deploy/debian/prerm
#!/bin/sh
systemctl stop prismarine
systemctl disable prismarine
